# backend/scripts/finance-app.service
[Unit]
Description=Finance App Spring Boot Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/finance-app
# Use which java to find Java path, fallback to common locations
ExecStart=/bin/bash -c 'JAVA_CMD=$(which java || echo /usr/bin/java); exec $JAVA_CMD $JAVA_OPTS -jar /opt/finance-app/finance-backend.jar'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=finance-app
# Give service time to start before considering it failed
TimeoutStartSec=60

# Load environment variables from file (if exists)
# Update this file with your actual values: /opt/finance-app/.env
EnvironmentFile=-/opt/finance-app/.env

# JVM Options (can be overridden in .env file)
Environment="JAVA_OPTS=-Xms512m -Xmx1024m -XX:+UseG1GC"

# Note: Set the following environment variables in /opt/finance-app/.env:
# DB_URL=jdbc:postgresql://your-rds-endpoint.region.rds.amazonaws.com:5432/financedb
# DB_USERNAME=postgres
# DB_PASSWORD=your_password
# ICICIDIRECT_APP_KEY=your_app_key
# ICICIDIRECT_CLIENT_SECRET=your_client_secret
# ICICIDIRECT_USER_ID=your_user_id
# ICICIDIRECT_PASSWORD=your_password
# ALLOWED_ORIGINS=http://localhost:3000,http://localhost:3001

[Install]
WantedBy=multi-user.target

